# src/main/deb/debian/postinst
#!/bin/sh
# postinst script for rock
#
# see: dh_installdeb(1)

set -e

. /usr/share/debconf/confmodule

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)

      # Create rock user if it doesn't exist.
      if ! id rock > /dev/null 2>&1 ; then
        adduser --system --home /var/lib/rock --no-create-home --disabled-password rock
      fi

      # rock file structure on Debian
      # /etc/rock: configuration
      # /usr/share/rock: executable
      # /var/lib/rock: data runtime
      # /var/log: logs

      rm -f /usr/share/rock
      ln -s /usr/share/rock-* /usr/share/rock

      if [ ! -e /var/lib/rock/conf ] ; then
        ln -s /etc/rock /var/lib/rock/conf
      fi

      # move installed r packages to the new library location
      rlibs=/var/lib/rock/R
      mkdir -p $rlibs
      if [ ! -e $rlibs/library ] ; then
        mkdir -p $rlibs/library
        if [ -e $rlibs/x* ] ; then
          rpkgs=`find $rlibs/x* -maxdepth 2 -mindepth 2 -type d`
          for pkg in $rpkgs
          do
            pkg_name=`basename $pkg`
            if [ ! -e $rlibs/library/$pkg_name ] ; then
              mv $pkg $rlibs/library
            fi
          done
        fi
      fi

      # Rserve package
      Rscript -e ".libPaths(.Library.site) ; if (!suppressWarnings(require(Rserve, quietly = TRUE))) { install.packages('Rserve', repos = 'https://www.rforge.net/') }"
      # utility packages
      Rscript -e ".libPaths(.Library.site) ; if (!suppressWarnings(require(unixtools, quietly = TRUE))) { install.packages('unixtools', repos = 'https://www.rforge.net/') }"
      Rscript -e ".libPaths(.Library.site) ; if (!suppressWarnings(require(remotes, quietly = TRUE))) { install.packages('remotes', repos = 'https://cloud.r-project.org/') }"
      Rscript -e ".libPaths(.Library.site) ; if (!suppressWarnings(require(BiocManager, quietly = TRUE))) { install.packages('BiocManager', repos = 'https://cloud.r-project.org/') }"
      # additional packages
      Rscript -e ".libPaths(.Library.site) ; if (!suppressWarnings(require(sqldf, quietly = TRUE))) { install.packages('sqldf', repos = 'https://cloud.r-project.org/') }"
      Rscript -e ".libPaths(.Library.site) ; if (!suppressWarnings(require(haven, quietly = TRUE))) { install.packages('haven', repos = 'https://cloud.r-project.org/') }"
      Rscript -e ".libPaths(.Library.site) ; if (!suppressWarnings(require(labelled, quietly = TRUE))) { install.packages('labelled', repos = 'https://cloud.r-project.org/') }"

      if [ ! -e /var/log/rock ] ; then
        mkdir /var/log/rock
        ln -s /var/log/rock /var/lib/rock/logs
      elif [ ! -e /var/lib/rock/logs ] ; then
        ln -s /var/log/rock /var/lib/rock/logs
      elif [ ! -L /var/lib/rock/logs ] ; then
        mv /var/lib/rock/logs/* /var/log/rock
        rmdir /var/lib/rock/logs
        ln -s /var/log/rock /var/lib/rock/logs
      fi

      chown -R rock:adm /var/lib/rock /var/log/rock /etc/rock
      chmod -R 750 /var/lib/rock /var/log/rock /etc/rock
      find /etc/rock/ -type f | xargs chmod 640

      if [ -e /etc/init.d/rock ]; then
        update-rc.d rock defaults >/dev/null
        rm /etc/init.d/rock
      fi
      systemctl daemon-reload
      systemctl enable rock
      systemctl start rock
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
